# recipes-core/project-arduinocli/files/arduinocli.service
[Unit]
Description=Arduino serial setup 
DefaultDependencies=no

[Service]
ExecStart=/bin/sh /usr/bin/cli.sh

[Install]
WantedBy=multi-user.target
